# preload/files/preload~
#!/sbin/runscript

description="Readahead daemon"

if [ -f /etc/conf.d/preload ]; then
	. /etc/conf.d/preload
fi

# Check for min memory
MIN_MEMORY=${MIN_MEMORY:-256}
free -m | awk  '/Mem:/ {exit ($2 >= ('"$MIN_MEMORY"'))?0:1}' || exit 0

# Check for ionice and use idle scheduling if available
RUNAS=""
IONICE="/usr/bin/ionice"
if [ -n "$IONICE_OPTS" ]; then
	if [ -x "$IONICE" ]; then
		RUNAS="$IONICE $IONICE_OPTS"
	else
		echo "ionice not found, running with normal io priority" >&2
	fi
fi

RETVAL=0

start() {
	ebegin "Starting preload"
	start-stop-daemon --start --exec "$RUNAS /usr/sbin/preload" -- $PRELOAD_OPTS
	eend $? "Failed to start preload"
}

stop() {
	ebegin "Stopping preload"
	start-stop-daemon --stop --exec /usr/sbin/preload
	eend $? "Failed to stop preload"
}
