# openrc/files/init.d/lvm
#!/sbin/runscript

depend(){
	before fsck localmount
	after mdadm
}

start() {
	ebegin "Initializing LVM support"
	vgscan --mknodes --ignorelockingfailure
	vgchange -ay --ignorelockingfailure
	eend $?
}

