# openrc/files/init.d/mdadm
#!/sbin/runscript

depend(){
	before fsck localmount
}

start() {
	ebegin "Initializing Linux SoftRaid support"
	mdadm -A -s
	eend 0
}
