# lirc-kernel-lts/files/lirc
#!/sbin/runscript

start() {
  ebegin "Starting lirc daemon: lircd and lircmd"
  start-stop-daemon --start --quiet --exec /usr/sbin/lircd -- --driver=devinput -d name=""
  start-stop-daemon --start --quiet --exec /usr/sbin/lircmd /etc/lirc/lircd.conf
  eend $?
}

stop() {
  ebegin "Stoping lirc daemon and lircmd"
  start-stop-daemon --stop --quiet --exec /usr/sbin/lircd
  start-stop-daemon --stop --quiet --exec /usr/sbin/lircmd
  eend $?
}