# teamviewer/files/teamviewer
#!/sbin/runscript

prog_path="/opt/teamviewer{ver}/tv_bin/teamviewerd"
pidfile="/var/run/teamviewerd.pid"

depend() {
   use logger dns
   need net xdm
   after xdm
}

start() {
  if ! [ -r ${pidfile} ] ; then
       ebegin "Starting teamviewerd"
       start-stop-daemon --start \
                         --quiet \
                         --pidfile ${pidfile} \
                         --make-pidfile \
                         --exec ${prog_path}
  fi
  eend $?
}

stop () {
  ebegin "Stopping teamviewerd"
  echo "Found pidfile ${pidfile} $(cat ${pidfile})"
  kill `cat ${pidfile}` >/dev/null 2>&1
  rm -rf ${pidfile} >/dev/null 2>&1
  eend $?
}

