# teamviewer/files/teamviewer
#!/sbin/runscript

prog="teamviewerd"
prog_path="/opt/teamviewer8/tv_bin/$prog"
prog_args="-d"
pidfile="/var/run/teamviewerd.pid"

depend() {
   use logger dns
   need net xdm
   after xdm
}

start() {
  if ! [ -r ${pidfile} ] ; then
       ebegin "Starting teamviewerd"
       start-stop-daemon --start \
                         --quiet \
                         --pidfile ${pidfile} \
                         --make-pidfile \
                         --exec ${prog_path}
  fi
  eend $?
}

stop () {
  ebegin "Stopping teamviewerd"
  if [ -r ${pidfile} ] ; then
       echo "Found pidfile ${pidfile} $(cat ${pidfile})"
       start-stop-daemon --stop \
                         --quiet \
                         --pidfile=${pidfile} --retry 20
  fi
  eend $?
}

