# x11vnc/files/x11vnc
#!/sbin/runscript
# -auth:
# gdm:  -auth /var/gdm/:0.Xauth
#       -auth /var/lib/gdm/:0.Xauth
# kdm:  -auth /var/lib/kdm/A:0-crWk72
#       -auth /var/run/xauth/A:0-crWk72
# xdm:  -auth /var/lib/xdm/authdir/authfiles/A:0-XQvaJk
#       dtlogin: -auth /var/dt/A:0-UgaaXa
# Sometimes the command "ps wwwwaux | grep auth" can reveal the file location.

PIDFILE=/var/run/x11vnc.pid

if [ "`ps -e | grep kdm`" ]; then
  AUTH='/var/lib/kdm/A:0-crWk72'
elif [ "`ps -e | grep gdm`" ]; then
  AUTH='/var/gdm/:0.Xauth'
elif [ "`ps -e | grep xdm`" ]; then
  AUTH='/var/lib/xdm/authdir/authfiles/A:0-XQvaJk'
fi

AUTH=`ps ax | grep X |awk '/-auth/{print $14;};'`

depend() {
   need xdm
   after xdm
}


start() {
  ebegin "Starting x11vnc"
  start-stop-daemon --start \
                    --quiet \
                    --pidfile ${PIDFILE} \
                    --make-pidfile \
                    --background \
                    --exec x11vnc -- \
                    -display :0 \
                    -auth ${AUTH} \
                    -xkb \
                    -rfbauth /etc/x11vnc/passwd \
                    -noncache \
                    -rfbport 5900 \
                    -o /var/log/x11vnc.log -loop
  if [ ! -f ${PIDFILE} ]; then
     eerror "Failed to start x11vnc"
     eend 1
     return 1
  fi
  eend $?
}

stop () {
  ebegin "Stopping x11vnc"
  if [ -r ${PIDFILE} ] ; then
       echo "Found pidfile $PIDFILE $(cat $PIDFILE)"
       start-stop-daemon --stop --quiet --pidfile=$PIDFILE --retry 20
  fi
  if [ 'ps -A | grep x11vnc' ] ; then
      echo "Killing all"
      killall -HUP x11vnc
  fi

  eend $?
}

