# synergy/files/synergy.initd
#!/sbin/runscript


depend() {
	need net
}

start() {
	emsg "Starting Synergy server"
	/usr/bin/synergys --config /etc/synergy.conf
	eend $?
}

stop() {
	emsg "Stopping Synergy server"
	killall synergys
	eend $?
}
