# x11vnc-misc/files/x11vnc
#!/sbin/runscript

name="x11vnc server daemon"

PIDFILE=/var/run/x11vnc.pid

depend() {
   use logger net
   need net
   after *
}

start() {
  if [ ! -f /etc/x11vnc/passwd ] ; then
     eerror "Start: sudo x11vnc --storepasswd /etc/x11vnc/passwd"
     return 1
  fi
  ebegin "Starting x11vnc"
  start-stop-daemon --start \
                    --quiet \
                    --pidfile ${PIDFILE} \
                    --make-pidfile \
                    --background \
                    --exec x11vnc -- \
                    -display :0 \
                    ${X11VNC_AUTH} \
                    -rfbauth /etc/x11vnc/passwd \
                    -rfbport ${X11VNC_PORT} \
                    ${X11VNC_OPTIONS} \
                    -o /var/log/x11vnc.log -loop
  if [ ! -f ${PIDFILE} ]; then
     eerror "Failed to start x11vnc"
     eend 1
     return 1
  fi
  eend $?
}

stop () {
  ebegin "Stopping x11vnc"
  if [ -r ${PIDFILE} ] ; then
       echo "Found pidfile ${PIDFILE} $(cat ${PIDFILE})"
       start-stop-daemon --stop --quiet x11vnc --retry 20
  fi
  eend $?
}

